# etc/cron.daily/linuxmuster-dehydrated
#!/bin/sh
/usr/sbin/linuxmuster-dehydrated -c >> /var/log/linuxmuster-dehydrated.log 2>&1
